# files/99-cd-audio-rip.rules
SUBSYSTEM=="block",KERNEL=="{{ cdrom_device_name }}",ACTION=="change",RUN+="/bin/systemctl start autocdripper"
